# bot/bot.service
[Unit]
Description=Rodep Bot with Auto-Upload
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/rodep/rodep-lore

# 起動時は startup.sh を叩く（ビルド・書き込み・upまでやってくれる）
# startup.sh の最後の "docker compose logs -f" は消しておいたほうが無難です
ExecStart=/bin/bash bot/startup.sh

# 停止時はしっかりコンテナを落とす
ExecStop=/usr/bin/docker compose down

User=rodep

[Install]
WantedBy=multi-user.target
